# apps/server/packaging/deb/prerm
#!/bin/sh
set -e

if [ "$1" = "remove" ]; then
  systemctl disable --now print-bridge.service || true
fi
